# main/deploy/summarizer-api.service
[Unit]
Description=Bedrock Summarizer API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=summarizer
Group=summarizer
WorkingDirectory=/opt/summarizer-api
EnvironmentFile=/etc/summarizer-api.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/summarizer-api/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8080 --workers 2
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
